# docs/systemd/scope-backend.service
[Unit]
Description=Scope Doc Backend
After=network.target

[Service]
User=scope
WorkingDirectory=/opt/scope_doc_gen/server
EnvironmentFile=/opt/scope_doc_gen/server/.env
ExecStart=/opt/scope_doc_gen/venv/bin/uvicorn server.api:app --host 127.0.0.1 --port 8010 --workers 2
Restart=always
RuntimeDirectory=scope-backend
RuntimeDirectoryMode=0755
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
