# deploy/scope-frontend.service
[Unit]
Description=Scope Doc Frontend (Next.js)
After=network.target

[Service]
Type=simple
User=scope
WorkingDirectory=/opt/scope_doc_gen/frontend
EnvironmentFile=/opt/scope_doc_gen/frontend/.env.local
ExecStart=/usr/bin/npm run start -- --hostname 127.0.0.1 --port 3300
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

